# root/etc/profile.d/10_term_default.sh
#!/bin/ksh
stty susp ''
stty intr ''
stty erase ''
stty kill ''
stty eof ''

umask 022
