# root/etc/profile.d/11_ko_lang.sh
#!/bin/ksh
stty cs8 -parenb -istrip -ixany

LANG=ko_KR.eucKR
export LANG
