# root/etc/profile.d/13_devenv.sh
#!/bin/ksh
export PATH=/usr/sbin:${PATH}
#export MANPATH=/opt/SUNWspro/man:/opt/SUNWspro/WS6U2/man:/usr/share/man
