# root/etc/profile.d/30_alitbase.sh
#!/bin/ksh
export ALTIBASE_HOME=/altibase/altibase_home

export PATH=${PATH}:${ALTIBASE_HOME}/bin
export LD_LIBRARY_PATH=${LD_LIBRARY_PATH}:${ALTIBASE_HOME}/lib
