# root/etc/profile.d/32_java.sh_backup_java6
#!/bin/ksh
export JAVA_HOME=/opt/java6
export PATH=${PATH}:${JAVA_HOME}/bin
