# root/etc/profile.d/32_java.sh
#!/bin/ksh
export JAVA_HOME=/opt/java1.5
export PATH=${PATH}:${JAVA_HOME}/bin
